# systemd/user/dropbox.service
[Unit]
Description=Dropbox Daemon

[Service]
ExecStart=/home/maurits/dotfiles/sync_dropbox.sh
Restart=always

[Install]
WantedBy=default.target

